# .config/systemd/user/autostart-sunshine.service
[Unit]
Description=Start Sunshine on session login

[Service]
ExecStart=%h/.local/bin/restart-sunshine

[Install]
WantedBy=graphical-session.target
