# .config/systemd/user/container@.service
# Run a container as a user
[Unit]
Description=On %u login, start container %i
After=network-online.target docker.service podman.service
Wants=network-online.target
RequiresMountsFor=%t/containers

[Service]
ExecStartPre=-/bin/sh -c 'podman stop -t 5 %i; podman rm -f %i'
ExecStart=/bin/sh -c '%h/.containers/bin/%i'
ExecStop=/bin/sh -c 'podman stop -t 5 %i; podman rm -f %i'

[Install]
WantedBy=default.target
