# .config/systemd/user/containerize@.service
[Unit]
Description=On %u login, start container %i

[Service]
Type=simple
ExecStart=%h/.containers/bin/%i

[Install]
WantedBy=basic.target
