# .config/systemd/user/frpc@.service
# Run fprc as a user
[Unit]
Description=On %u login, start frpc %i
After=network-online.target syslog.target remote-fs.target nss-lookup.target
Wants=network-online.target

[Service]
ExecStart=%h/.local/bin/jail frpc -- -c %i.toml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
