# .config/systemd/user/healthcheck@.service
[Unit]
Description=Healthcheck ping @ healthcheck.io %i

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -s 'https://hc-ping.com/%i'
