# dots/scripts/linux/pyamsoft-settings/hosts/skip:aws-/etc/systemd/system/container@scrutiny.service.d/override.conf
[Service]
# Does not require networking
Environment=PYAM_CHECK_NETWORK=0
