# dots/system/etc/systemd/system/NetworkManager.service.d/override.conf
[Unit]
Before=libvirtd.service mullvad-vpn.service tailscaled.service
