# dots/scripts/linux/pyamsoft-settings/hosts/skip:aws-/etc/systemd/user/container@beszel-hub.service.d/override.conf
[Service]
# No need to check networking
Environment=PYAM_CHECK_NETWORK=0
