# .config/systemd/user/container@blocky.service.d/override.conf
[Service]
# Clear old
ExecStartPre=

# Check the file exists, do not attempt to start otherwise
ExecStartPre=/bin/test -x %h/.containers/bin/%i

# Check the network is online before starting (wait up to 30 seconds)
# Blocky IS our DNS, so not check for DNS here.
ExecStartPre=/bin/sh -c '[ -x /opt/pyamsoft/bin/im-online ] && /opt/pyamsoft/bin/im-online'

# Teardown before starting, if this 'fails' because nothing else is running, we don't mind
ExecStartPre=-/opt/pyamsoft/bin/container-teardown %i
