# dots/scripts/linux/pyamsoft-settings/hosts/skip:vps--/etc/systemd/user/org.gnome.Shell@wayland.service.d/override.conf
# Disable hardware cursor for Moonlight/Sunshine to render cursor correctly
[Service]
Environment=MUTTER_DEBUG_DISABLE_HW_CURSORS=1
