# dots/scripts/linux/pyamsoft-settings/hosts/nerd2/etc/systemd/system/prepare-mergerfs-MassMedia.service
[Unit]
Description=Prepare mergerfs for /mnt/MassMedia/

[Service]
Type=notify
RemainAfterExit=yes
ExecStart=/opt/pyamsoft/sbin/prepare-mergerfs-MassMedia

[Install]
WantedBy=default.target
