# dots/scripts/linux/pyamsoft-settings/hosts/aws-/etc/systemd/system/crowdsec-firewall-bouncer.service.d/override.conf
[Unit]
Description=The firewall bouncer for CrowdSec
After=container-crowdsec.service container@crowdsec.service
