# .config/bash/profile.d/01-env
# shellcheck shell=bash

# Source bash_profile if the environment is not setup
if [ -z "${PYAMSOFT_ENVIRONMENT}" ]; then
  # shellcheck disable=SC1091
  [ -f "${HOME}"/.environment ] && . "${HOME}"/.environment
fi

return 0
