# .config/bash/profile.d/01-umask.shell
# shellcheck shell=bash

# Strict umask
if command -v umask >/dev/null; then
  umask 077 || return 1
fi

return 0

# vim: set filetype=sh:
