# .config/bash/profile.d/01-umask.shell
# shellcheck shell=bash

# Strict umask
if command -v umask >/dev/null; then
  _current_umask="$(umask)"

  case "${_current_umask}" in
  077 | 0077)
    ;;
  *)
    umask 077 || {
      unset _current_umask
      return 1
    }
    ;;
  esac

  unset _current_umask
fi

return 0

# vim: set filetype=sh:
