# .config/bash/profile.d/03-completion
# shellcheck shell=bash

# Move bash_completion file
export BASH_COMPLETION_USER_FILE="${XDG_CONFIG_HOME}/bash/completions"
