# .config/bash/profile.d/03-history.env
# shellcheck shell=bash

case "$-" in
*i*)
  # Continue as normal
  ;;
*)
  # Stop if not running interactively
  return 0
  ;;
esac

# Disable the bash_history file
export HISTCONTROL="ignoreboth"
export HISTFILE="${XDG_DATA_HOME}/bash/history"

return 0

# vim: set filetype=sh:
