# .config/bash/profile.d/03-history.env
# shellcheck shell=bash

# Disable the bash_history file
export HISTCONTROL="ignoreboth"
export HISTFILE="${XDG_CONFIG_HOME}/bash/history"

return 0

# vim: set filetype=sh:
