# .config/bash/profile.d/05-pyamsoft-path.env
# shellcheck shell=bash

if [ -d "/opt/pyamsoft/" ]; then
  # If we are on a system with /opt
  PYAMSOFT_PATH="/opt/pyamsoft"
else
  PYAMSOFT_PATH="${HOME}/.local/pyamsoft"
fi

# Export the PATH
export PYAMSOFT_PATH

return 0

# vim: set filetype=sh:
