# .config/bash/profile.d/06-editor.env.macos
# shellcheck shell=bash

if [ -n "${EDITOR}" ]; then
  return 0
fi

# User editor can be from a local script
# Use absolute path to editor
_editors="/opt/homebrew/bin/nvim"
for e in ${_editors}; do
  if [ -f "${e}" ]; then
    export EDITOR="${e}"
    export SYSTEMD_EDITOR="${e}"
    export SUDO_EDITOR="${e}"
    export DIFFPROG="${e} -d"
    break
  fi
done
unset _editors

return 0

# vim: set filetype=sh:
