# .config/bash/profile.d/06-pyamsoft-bin.env
# shellcheck shell=bash

# pyamsoft scripts
if [ -d "${PYAMSOFT_PATH}/bin" ]; then
  prepend_path "${PYAMSOFT_PATH}/bin"
fi

if [ -d "${PYAMSOFT_PATH}/sbin" ]; then
  prepend_path "${PYAMSOFT_PATH}/sbin"
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
