# .config/bash/profile.d/06-pyamsoft-opt.env
# shellcheck shell=bash

# pyamsoft scripts
if [ -d "/opt/pyamsoft/bin" ]; then
  prepend_path "/opt/pyamsoft/bin"
fi

if [ -d "/opt/pyamsoft/sbin" ]; then
  prepend_path "/opt/pyamsoft/sbin"
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
