# .config/bash/profile.d/20-audio.env.linux
# shellcheck shell=bash

# Help with Pulse audio crackling in some programs
export PULSE_LATENCY_MSEC=60

return 0

# vim: set filetype=sh:
