# .config/bash/profile.d/20-docker.env.macos
# shellcheck shell=bash

# Docker scripts
if [ -d "${HOME}/.docker/bin" ]; then
  PATH="${HOME}/.docker/bin:${PATH}"
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
