# .config/bash/profile.d/20-guake.env.linux
# shellcheck shell=bash

# Run Guake under wayland
# https://github.com/Guake/guake/pull/1975
export GUAKE_ENABLE_WAYLAND=1

return 0

# vim: set filetype=sh:
