# .config/bash/profile.d/20-gnupg.env
# shellcheck shell=bash

if [ -z "${GNUPGHOME}" ]; then
  # Suggestions from xdg-ninja
  export GNUPGHOME="${XDG_DATA_HOME}/gnupg"
  mkdir -p "${GNUPGHOME}"
fi

return 0

# vim: set filetype=sh:
