# .config/bash/profile.d/21-docker-socket.env.macos
# shellcheck shell=bash

# Docker socket
if [ -d "${HOME}/.docker/run" ] && [ -S "${HOME}/.docker/run/docker.sock" ]; then
  if [ -z "${DOCKER_HOST}" ]; then
    DOCKER_HOST="unix://${HOME}/.docker/run/docker.sock"
    export DOCKER_HOST
  fi
fi

return 0

# vim: set filetype=sh:
