# .config/bash/profile.d/22-qt-style.env.linux
# shellcheck shell=bash

# For QT applications running on GNOME, instead of fusion force the visual style to match the GNOME desktop
# This SHOUlD be pulled in by a user level flatpak environment too
#
# Just assume Adwaita
if [ -n "${XDG_CURRENT_DESKTOP}" ] && [ "${XDG_CURRENT_DESKTOP}" = "GNOME" ]; then
  export QT_STYLE_OVERRIDE=Adwaita-Dark
  export QT_QPA_PLATFORMTHEME=gtk3
fi

return 0

# vim: set filetype=sh:
