# .config/bash/profile.d/30-jail.env.linux
# shellcheck shell=bash

# Jail scrips
if command -v jail >/dev/null; then
  if [ -d "${HOME}/.jail/bin" ]; then
    PATH="${HOME}/.jail/bin:${PATH}"
  fi
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
