# .config/bash/profile.d/50-containers.env
# shellcheck shell=bash

# Docker scrips
if command -v docker >/dev/null || command -v podman >/dev/null; then
  if [ -d "/opt/pyamsoft/containers/bin" ]; then
    PATH="/opt/pyamsoft/containers/bin:${PATH}"
  fi
fi

# Distrobox scrips
if command -v distrobox >/dev/null; then
  if [ -d "/opt/pyamsoft/distrobox/bin" ]; then
    PATH="/opt/pyamsoft/distrobox/bin:${PATH}"
  fi
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
