# .config/bash/profile.d/50-flatpak.env.linux
# shellcheck shell=bash

# Flatpak scrips
if command -v flatpak >/dev/null; then
  if [ -d "/opt/pyamsoft/flatpak/bin" ]; then
    prepend_path "/opt/pyamsoft/flatpak/bin"
  fi
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
