# .config/bash/profile.d/50-jail.env.linux
# shellcheck shell=bash

# Jail scrips
if command -v jail >/dev/null; then
  if [ -d "/opt/pyamsoft/jail/bin" ]; then
    prepend_path "/opt/pyamsoft/jail/bin"
  fi
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
