# .config/bash/profile.d/50-pyamsoft-flatpak.env.linux
# shellcheck shell=bash

# Flatpak scrips
if command -v flatpak >/dev/null; then
  if [ -d "${PYAMSOFT_PATH}/flatpak/bin" ]; then
    prepend_path "${PYAMSOFT_PATH}/flatpak/bin"
  fi
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
