# .config/bash/profile.d/50-pyamsoft-jail.env.linux
# shellcheck shell=bash

# Jail scrips
if command -v jail >/dev/null; then
  if [ -d "${PYAMSOFT_PATH}/jail/bin" ]; then
    prepend_path "${PYAMSOFT_PATH}/jail/bin"
  fi
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
