# .config/bash/profile.d/52-dircolors
# shellcheck shell=bash

if [ -z "${LS_COLORS}" ]; then
  if command -v dircolors >/dev/null; then
    eval "$(dircolors -b)" || return 1
  fi
fi

return 0
