# .config/bash/profile.d/60-containers.env
# shellcheck shell=bash

# Docker scrips
if command -v docker >/dev/null || command -v podman >/dev/null; then
  if [ -d "${PYAMSOFT_PATH}/containers/bin" ]; then
    prepend_path "${PYAMSOFT_PATH}/containers/bin"
  fi
fi

# Distrobox scrips
if command -v distrobox >/dev/null; then
  if [ -d "${PYAMSOFT_PATH}/distrobox/bin" ]; then
    prepend_path "${PYAMSOFT_PATH}/distrobox/bin"
  fi
fi

# Export the PATH
export PATH

return 0

# vim: set filetype=sh:
