# dots/scripts/linux/pyamsoft-settings/hosts/skip:aws-/etc/systemd/system/unbound-trusted-keys@.timer
[Unit]
Description=Run trusted-keys monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
