# dots/scripts/linux/pyamsoft-settings/hosts/vps--/etc/systemd/user/container@beszel-remote.service.d/override.conf
[Service]
# We are running in remote mode
Environment=BESZEL_MODE=agent

# Run the beszel script
Environment=PYAM_CONTAINER=/opt/pyamsoft/bin/beszel
