# dots/scripts/containers/scrutiny/scrutiny-collector-metrics.timer
[Unit]
Description=Scrutiny Collector Metrics Daily Timer

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
