# .config/systemd/user/org.gnome.Shell@wayland.service.d/override.conf
[Service]
Environment=MUTTER_DEBUG_DISABLE_HW_CURSORS=1
