# dots/scripts/linux/pyamsoft-settings/hosts/all/etc/systemd/user/container@beszel-hub.service.d/override.conf
[Service]
# We are running in hub mode
Environment=BESZEL_MODE=hub

# Run the beszel script
Environment=PYAM_CONTAINER=/opt/pyamsoft/bin/beszel
