# dots/scripts/linux/pyamsoft-settings/hosts/nerd2/etc/systemd/system/container@crowdsec.service.d/override.conf
[Service]
Environment=CROWDSEC_MODE=lapi
