# dots/scripts/linux/pyamsoft-settings/hosts/nerd1/etc/systemd/user/cameractrlsd.service
[Unit]
Description=CameraCtrls daemon - restore control values

[Service]
Type=simple
ExecStart=/usr/bin/flatpak --user run --command='cameractrlsd.py' hu.irl.cameractrls

[Install]
WantedBy=graphical-session.target
