# dots/scripts/linux/pyamsoft-settings/hosts/skip:aws-/etc/systemd/user/container@blocky.service.d/override.conf
[Service]
# Blocky IS our DNS, so not check for DNS here.
Environment=PYAM_CHECK_NETWORK_DNS=0
