# dots/scripts/linux/pyamsoft-settings/hosts/vps--/etc/systemd/system/container@crowdsec.service.d/override.conf
[Service]
Environment=CROWDSEC_MODE=agent
