# dots/scripts/linux/pyamsoft-settings/hosts/skip:aws-/etc/systemd/system/unbound.service.d/override.conf
[Service]
ExecPreStart=/usr/bin/chown unbound:unbound -R /etc/unbound/rpz.d
