# dots/scripts/linux/pyamsoft-settings/hosts/skip:aws-/etc/systemd/system/unbound.service.d/override.conf
[Unit]
After=unbound-prepare.service
Requires=unbound-prepare.service
