# dots/scripts/linux/pyamsoft-settings/hosts/nerd2/etc/systemd/user/container-freegames-claim.timer
[Unit]
Description=Claim freegames using vogler/free-games-claimer

[Timer]
OnCalendar=0/4:00:00
Persistent=true
Unit=container@freegames.service

[Install]
WantedBy=timers.target
