# system/etc/systemd/system/adjust-ryzen.service
[Unit]
Description=Adjust Ryzen CPU
After=default.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/local/bin/adjust-ryzen
ExecStop=/usr/local/bin/adjust-ryzen

[Install]
WantedBy=default.target
